# couchdb_from_official_on_debian/roles/couchdb_from_official_on_debian/templates/couchdb.service.j2
[Unit]
Description=CouchDB service
After=network.target

[Service]
Type=simple
User=couchdb
ExecStart=/opt/couchdb/bin/couchdb -o /dev/null -e /dev/null

[Install]
WantedBy=multi-user.target
